# packer/scripts/crucible-appliance-startup.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/crucible-appliance-startup.sh

[Install]
WantedBy=default.target